# linux_service/etc.systemd.system/code-server-eth.service
[Unit]
Description=Code Server

[Service]
ExecStart=/home/ubuntu/code-server/start-eth.sh

[Install]
WantedBy=multi-user.target
