# linux_service/etc.systemd.system/code-server-wlan.service
[Unit]
Description=Code Server

[Service]
ExecStart=/home/ubuntu/code-server/start-wlan.sh

[Install]
WantedBy=multi-user.target
