# linux_service/etc.systemd.system/config-wlan.service
[Unit]
Description=Configuration wlan0

[Service]
ExecStart=/home/ubuntu/web-menu/config_wlan.sh

[Install]
WantedBy=network-online.target
