# linux_service/etc.systemd.system/eth-butterfly.service
[Unit]
Description=Butterfly Terminal Server

[Service]
ExecStart=/usr/local/bin/eth-butterfly.py

[Install]
WantedBy=multi-user.target
