# linux_service/etc.systemd.system/pioneer-bricks-wlan.service
[Unit]
Description=Geoscan Pioneer Bricks

[Service]
ExecStart=/home/ubuntu/pioneer-bricks/start-wlan.sh

[Install]
WantedBy=multi-user.target
