# linux_service/etc.systemd.system/web-wlan-menu.service
[Unit]
Description=Pioneer Web-Menu

[Service]
ExecStart=/home/ubuntu/web-menu/start-wlan.sh

[Install]
WantedBy=multi-user.target
