# linux_service/etc.systemd.system/web-eth-menu.service
[Unit]
Description=Pioneer Web-Menu

[Service]
ExecStart=/home/ubuntu/web-menu/start-eth.sh

[Install]
WantedBy=multi-user.target
