# linux_service/etc.systemd.system/wlan-butterfly.service
[Unit]
Description=Butterfly Terminal Server

[Service]
ExecStart=/usr/local/bin/wlan-butterfly.py

[Install]
WantedBy=multi-user.target
